# systemd/.config/systemd/user/rclone-sync.timer
[Unit]
Description=Rclone Bisync Timer

[Timer]
OnBootSec=5min
OnUnitActiveSec=60min

[Install]
WantedBy=timers.target