# systemd/.config/systemd/user/void-auditor.timer
[Unit]
Description=Run Void Auditor Daily

[Timer]
# Triggers once a day at midnight
OnCalendar=daily
# Ensures it runs immediately on boot if the laptop was off at midnight
Persistent=true

[Install]
WantedBy=timers.target
