# systemd/.config/systemd/user/void-auditor.service
[Unit]
Description=Native Void: Rclone Silent Auditor

[Service]
Type=oneshot
ExecStart=/home/manokel/.local/bin/rclone_auditor.sh
